# filamama-backend.service
[Unit]
Description=FilaMama Backend - FastAPI File Manager Service
After=network.target

[Service]
Type=simple
User=flatstone
Group=flatstone
WorkingDirectory=/home/flatstone/FLATSTONE/FilaMama/backend
Environment="PATH=/home/flatstone/FLATSTONE/FilaMama/backend/venv/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin"
ExecStart=/home/flatstone/FLATSTONE/FilaMama/backend/venv/bin/uvicorn app.main:app --host 0.0.0.0 --port 8011
Restart=always
RestartSec=10

# Security hardening
# NOTE: ProtectHome is intentionally NOT set — root_path is under /home/flatstone, which
# the service must read/write. ProtectSystem=full keeps /usr, /boot and /etc read-only.
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictSUIDSGID=true
RestrictRealtime=true
RestrictNamespaces=true
LockPersonality=true

[Install]
WantedBy=multi-user.target
